# default_files/etc/skel/.config/autostart/beallit.desktop
[Desktop Entry]
Name=Beallitó
Exec=beallit
Icon=
Type=Application
Categories=System;Utility;

