# default_files/etc/skel/.config/autostart/telepit.desktop
[Desktop Entry]
Name=Telepítő
Exec=Magyarcs
Icon=/usr/share/icons/Faenza-Dark/places/22/start-here.png
Type=Application
Categories=System;Utility;


